# runner-resource-caps.service
[Unit]
Description=runner-resource-caps

[Service]
ExecStart=/root/runner-resource-caps.sh
Restart=always

[Install]
WantedBy=multi-user.target
